# deploy/systemd/dbaops-a2a.service
# DBAOps native A2A 서버 (:9102) — DBAOps RCA 에이전트를 A2A로 직접 노출.
# invoke_single 을 직접 호출하므로 mcp-router(:9000) + Bedrock 필요 (dbaops-agent와 동일).
[Unit]
Description=DBAOps native A2A server
After=network-online.target dbaops-mcp-router.service
Wants=network-online.target
Requires=dbaops-mcp-router.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DBAOPS__
EnvironmentFile=/etc/dbaops/dbaops.env
Environment=PYTHONPATH=__DBAOPS__
Environment=GATEWAY_ENDPOINT=http://127.0.0.1:9000/mcp
Environment=DBAOPS_DATA_DIR=__DATA__
Environment=DBAOPS_A2A_HOST=0.0.0.0
Environment=DBAOPS_A2A_PORT=9102
Environment=DBAOPS_A2A_URL=http://127.0.0.1:9102/
Environment=PERF_A2A_URL=http://127.0.0.1:9100
ExecStart=__VENV__/bin/python a2a_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
