# deploy/systemd/dbperf-ops-facade.service
# DBAOps A2A 파사드 (:9101) — DBAOps agent(:8080 HTTP)를 A2A로 노출. 양방향.
[Unit]
Description=DBPerf ops-facade (DBAOps agent A2A bridge)
After=network-online.target dbaops-agent.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PERF__
EnvironmentFile=/etc/dbaops/dbaops.env
Environment=A2A_OPS_HOST=0.0.0.0
Environment=A2A_OPS_PORT=9101
Environment=A2A_OPS_URL=http://127.0.0.1:9101
Environment=DBAOPS_AGENT_URL=http://127.0.0.1:8080/invocations
Environment=PERF_A2A_URL=http://127.0.0.1:9100
ExecStart=__VENV__/bin/python a2a_ops_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
