# dbaops/deploy/ec2-vanilla/systemd/dbaops-agent.service
# DBAOps Agent — 프로세스 1개 · 포트 1개(:8080):
#   A2A (perf 에이전트) + /invocations NDJSON 스트리밍 (UI·Slack) 통합 서빙
[Unit]
Description=DBAOps Agent
After=network-online.target dbaops-mcp-router.service
Wants=network-online.target
Requires=dbaops-mcp-router.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO__
EnvironmentFile=/etc/dbaops/dbaops.env
Environment=PYTHONPATH=__REPO__
Environment=GATEWAY_ENDPOINT=http://127.0.0.1:9000/mcp
Environment=DBAOPS_DATA_DIR=__DATA__
Environment=PORT=8080
Environment=DBAOPS_A2A_URL=http://127.0.0.1:8080/
Environment=PERF_A2A_URL=http://127.0.0.1:9100
ExecStart=__VENV__/bin/python a2a_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
